# timers/lookingback.timer
[Unit]
Description=Run lookingback.sh every day at 12:00

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true

[Install]
WantedBy=timers.target
