# timers/pdnssoc.timer
[Unit]
Description=Run pdnssoc.rb every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
